# infra/ansible/roles/edge_robot/templates/freedroid-health.service.j2
[Unit]
Description=Free-Droid vital-function health check (network / hardware / software)
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=oneshot
# Runs as root: it may `systemctl restart` wg-quick@wg0 / ollama / freedroid
# (self-heal) and writes /run/freedroid/health.json + the safe-mode flag.
# Assert we're on the Pi so a transient device-tree read can't silently SKIP the
# CRITICAL hardware checks (a missing device then FAILs -> safe-mode, fail-safe).
Environment=FREEDROID_ASSUME_PI=1
ExecStart={{ robot_venv_bin }}/freedroid-health
